# supervisor.conf
[program:entrant_gunicorn]
command=/home/entrant/dev/flask_api/start.sh
user=entrant
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true